# deploy/ubuntu/canvas-video.service
[Unit]
Description=Private SJTU Canvas Video Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=canvas-video
Group=canvas-video
WorkingDirectory=/opt/canvas-video/current
Environment=SJTU_CANVAS_CONFIG=/etc/canvas-video/config.toml
Environment=SJTU_REAL_PROTOCOL_TEST=1
Environment=RUST_LOG=server=info,canvas_core=info
ExecStart=/opt/canvas-video/current/bin/sjtu-canvas-video-server
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
KillSignal=SIGTERM
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
LimitNOFILE=65536
TasksMax=256

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
